# docker/php-fpm/crontab
* * * * * cd /application && /usr/local/bin/php artisan schedule:run >> storage/logs/crontab.out 2>&1
